# roles/gardener/templates/sensor.service.j2
[Unit]
Description=Sensor reader

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sensor -o /run/sensor/prometheus
