# roles/gardener/templates/sensor.service.j2
[Unit]
Description=Sensor reader

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sensor /etc/sensor.conf /var/run/sensor/metrics
