# roles/gardener/templates/sensor.timer.j2
[Unit]
Description=Sensor reader

[Timer]
OnStartupSec=10
OnUnitActiveSec=1

[Install]
WantedBy=timers.target
