# roles/gardener/templates/sensor.timer.j2
[Unit]
Description=Sensor reader

[Timer]
OnActiveSec=1

[Install]
WantedBy=timers.target
