# roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus

[Service]
ExecStart=/usr/local/bin/prometheus --config.file=/etc/prometheus.yml
Restart=always

[Install]
WantedBy=multi-user.target
